# serverside/displayIP.conf
[program:getIPAddress]
command = python app.py
directory = /home/pi/raspberryip
autostart = true
autorestart = true